# os/os/package/drafter-agent/drafter-agent.service.in
[Unit]
Description=Drafter agent daemon
After=network.target @DRAFTER_AGENT_SYSTEMD_DEPENDENCY@
Requires=@DRAFTER_AGENT_SYSTEMD_DEPENDENCY@

[Service]
Type=simple
ExecStart=/usr/bin/drafter-agent --vsock-port @DRAFTER_AGENT_VSOCK_PORT@ --before-suspend-cmd @DRAFTER_AGENT_BEFORE_SUSPEND_CMD@ --after-resume-cmd @DRAFTER_AGENT_AFTER_RESUME_CMD@
StandardOutput=journal+console
StandardError=journal+console
Restart=always

[Install]
WantedBy=multi-user.target
