# os/package/drafter-liveness/drafter-liveness.service.in
[Unit]
Description=Drafter liveness daemon
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/drafter-liveness --vsock-port @DRAFTER_LIVENESS_VSOCK_PORT@
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
